# raevtar.service.tmpl
[Unit]
Description=Raevtar — personal blog and dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{RAEVTAR_USER}}
WorkingDirectory={{RAEVTAR_HOME}}
ExecStart={{RAEVTAR_HOME}}/raevtar
EnvironmentFile={{RAEVTAR_HOME}}/.env.production
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
